# AI-teams-controller-public/scripts/ai-teams-frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-teams.frontend</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/pnpm</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sonph36/dev/tools/AI-teams-controller/frontend</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>3334</string>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/ai-teams-frontend.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ai-teams-frontend.err</string>
</dict>
</plist>
